# services/secrets-api/hermes-secrets-api.service
[Unit]
Description=Hermes Secrets API
Documentation=https://github.com/ruimgoncalves/telegram-hermes-miniapp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/telegram-hermes-miniapp/services/secrets-api
Environment=PYTHONUNBUFFERED=1
Environment=HERMES_HOME=/root/.hermes
Environment=SECRETS_API_HOST=127.0.0.1
Environment=SECRETS_API_PORT=9876
ExecStart=/root/projects/telegram-hermes-miniapp/services/secrets-api/.venv/bin/python server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/.hermes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
